# passenger/supervisor.conf
[program:passenger]
command=passenger start
directory=/www
autorestart=true
