# softether/supervisord.conf
[supervisord]
nodaemon=true

[program:vpnserver]
command=bash -c 'cp /opt/vpnserver/* /usr/local/vpnserver/; exec /usr/local/vpnserver/vpnserver execsvc'
