# transmission/crontab
# m h  dom mon dow   command
* * * * * prss `cat ~/watch/prss-key` ~/watch >> /dev/null 2>&1
